# bootstrap/supervisord.conf
[program:randomly]
user=randomly
command=/srv/app/randomly/venv/bin/gunicorn -b 0.0.0.0:8080 randomly.wsgi
directory=/srv/app/randomly
autostart=true
autorestart=true
stderr_logfile=/var/log/randomly/randomly.err.log
stdout_logfile=/var/log/randomly/randomly.out.log
